# docker/images/rails/supervisord.conf
[supervisord]
nodaemon=true

[program:unicorn]
command=/start-server.sh

[program:nerve]
command=/start-nerve.sh