# docker/images/synapse/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:synapse]
command=synapse -c /synapse-config.json

[program:haproxy]
command=haproxy -f /etc/haproxy/haproxy.cfg